# archived/salt/administration/files/aws_im.conf
################################################################
# aws_im Manager Service

description     "This upstart configuration takes care of starting up the aws_im manager"
author          "william"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

# automatically restart service if it stopped unexpectedly
respawn

#expect fork

# how the script is actually started, taking care to start it 
script
  exec sudo  /usr/bin/python /srv/salt/keymanager/aws_im.py
end script

